# etc/systemd/system/zfs-scrub@.service.in
[Unit]
Description=zpool scrub on %i
Documentation=man:zpool-scrub(8)
DefaultDependencies=no
Requires=zfs.target
After=zfs.target
ConditionACPower=true
ConditionPathIsDirectory=/sys/module/zfs

[Service]
Type=oneshot
ExecStart=@sbindir@/zpool scrub %i
