# etc/systemd/system/zfs-scrub@.service.in
[Unit]
Description=zpool scrub on %i
Documentation=man:zpool-scrub(8)
Requires=zfs.target
After=zfs.target
ConditionACPower=true
ConditionPathIsDirectory=/sys/module/zfs

[Service]
KillSignal=SIGINT
ExecStart=@sbindir@/zpool scrub %i
